# scripts/agents/systemd/faktomat-backup.service
[Unit]
Description=Faktomat Postgres backup + Discord report
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ewwe
Group=ewwe
WorkingDirectory=/home/ewwe/playspace/random
ExecStart=/home/ewwe/playspace/random/scripts/agents/backup_verifier.sh
SuccessExitStatus=0 1

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
